# systemd/system/disable-wakeup.service
[Unit]
Description=Disable HXCi USB wakeup

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'echo XHC > /proc/acpi/wakeup'

[Install]
WantedBy=default.target
