# systemd/backup-notes.timer
[Unit]
Description=Backup notes daily

[Timer]
Persistent=true  
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=default.target