# systemd/user/backup.timer
[Unit]
Description=Backup

[Timer]
Persistent=true  
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=default.target