# systemd/user/backup-notes.service
[Unit]
Description=Backup my notes

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/home/samim/.bin/backup-notes'