# systemd/user/feh-wallpaper.timer
[Unit]
Description=Random wallpaper with feh

[Timer]
OnUnitActiveSec=300s
Unit=feh-wallpaper.service

[Install]
WantedBy=default.target