# systemd/user/prismatik.service
[Unit]
Description=Prismatik

[Service]
ExecStart=/usr/bin/prismatik --nogui

[Install]
WantedBy=default.target
