# systemd/user/sync-agenda-svg.service
[Unit]
Description=Sync agenda svg

[Service]
Type=oneshot
ExecStart=/bin/bash -c /home/samim/.bin/sync-agenda-svg
