# systemd/user/sync-cal.service
[Unit]
Description=Sync cal

[Service]
Type=oneshot
ExecStart=/bin/bash -c /home/samim/.bin/sync-cal
