# crontab/root.crontab
# vim:set ft=crontab:
0 6 * * * shutdown now
1 6 * * * poweroff
